# deploy/systemd/mywebapp-container.service
[Unit]
Description=mywebapp container stack (Лабораторна №3 — CI/CD deploy)
After=docker.service network-online.target
Wants=docker.service network-online.target
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/mywebapp
EnvironmentFile=/opt/mywebapp/.env
ExecStartPre=/usr/bin/docker compose -f /opt/mywebapp/compose.prod.yml pull
ExecStart=/usr/bin/docker compose -f /opt/mywebapp/compose.prod.yml up -d --remove-orphans
ExecStop=/usr/bin/docker compose -f /opt/mywebapp/compose.prod.yml down
Restart=on-failure
RestartSec=5s
TimeoutStartSec=180

[Install]
WantedBy=multi-user.target
